# docker/supervisord_app.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:backend]
command=python3 /app/engineering_service.py
autostart=true
autorestart=true
stdout_logfile=/tmp/logs/backend.log
stderr_logfile=/tmp/logs/backend.err
environment=PYTHONUNBUFFERED="1",ENGINEERING_SERVICE_PORT="8000",ENGINEERING_SERVICE_HOST="0.0.0.0"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/tmp/logs/nginx.log
stderr_logfile=/tmp/logs/nginx.err
